# files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyGS0::respawn:/sbin/getty -L -w -n 115200 /dev/ttyGS0
